# cron/crontab
0 1 * * * /opt/venv/bin/python3 /birthday_update_grpc_script.py > /proc/1/fd/1 2>/proc/1/fd/2
