# cron/crontab
0 1 * * * python /birthday_update_grpc_script.py 